# incubation/setmy-info-scripts/resources/etc/profile.d/example.sh
export SETMY_INFO_PATH=/opt/setmy.info
export PATH=${SETMY_INFO_PATH}/bin:$PATH
export MANPATH=${SETMY_INFO_PATH}/man:${MANPATH}
