# setmy-info-scripts/resources/etc/profile.d/setmy-info.sh
export SMI_PROVIDER=setmy.info
export SETMY_INFO_PATH=/opt/${SMI_PROVIDER}
export PATH=${SETMY_INFO_PATH}/bin:$PATH
export MANPATH=${SETMY_INFO_PATH}/man:${MANPATH}
